# supervisord.conf
[supervisord]
nodaemon=true

[program:sshd]
command=/usr/sbin/sshd -D
autorestart=true

[program:jenkins]
command=java -server -jar jenkins.war
autorestart=true

